# build/laravel-s.supervisor.conf
[program:laravel-s]
process_name=%(program_name)s
command=php /app/bin/laravels start -i
autostart=true
autorestart=true
startretries=3
user=root
numprocs=1
stdout_logfile=/app/storage/logs/%(program_name)s.log
stdout_logfile_maxbytes=0
stderr_logfile=/app/storage/logs/%(program_name)s_error.log
stderr_logfile_maxbytes=0
